# script/dev.harness.reap-docker-disk.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>dev.harness.reap-docker-disk</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <!-- Installed WITHOUT --apply on purpose: soak in dry-run, read the log, then add
         "--apply --prune-nested" here and re-run install-reaper-launchd.sh. -->
    <string>exec "__REPO__/script/reap-docker-disk.sh"</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <!-- launchd fires a missed calendar interval on the next wake; cron silently skips it,
       and this is a laptop that is asleep at 04:00 most nights. -->
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
  <key>RunAtLoad</key><false/>
  <key>StandardOutPath</key><string>/tmp/harness-reap.out</string>
  <key>StandardErrorPath</key><string>/tmp/harness-reap.err</string>
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityIO</key><true/>
  <key>Nice</key><integer>10</integer>
</dict>
</plist>
